# patch/systemd/node-agent.service
[Unit]
Description=ClusterOS node-agent
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'if [ -f /root/patch/apply-patch.sh ] && [ ! -f /var/lib/clusteros/.bootstrapped ]; then bash /root/patch/apply-patch.sh >> /var/log/clusteros-boot.log 2>&1 && mkdir -p /var/lib/clusteros && touch /var/lib/clusteros/.bootstrapped; fi'
ExecStart=/usr/local/bin/node-agent --log-level info start --foreground
TimeoutStartSec=0
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
EnvironmentFile=-/etc/cluster-os/node-agent.env

[Install]
WantedBy=multi-user.target