# images/ubuntu/systemd/node-agent.service
[Unit]
Description=Cluster OS Node Agent
Documentation=https://github.com/cluster-os
After=network-online.target
Wants=network-online.target
Before=slurm-wlm.service k3s.service

[Service]
Type=simple
User=root
Group=root

# Node agent binary
ExecStart=/usr/local/bin/node-agent

# Environment
Environment="NODE_CONFIG_PATH=/etc/cluster-os/node.yaml"
Environment="NODE_DATA_DIR=/var/lib/cluster-os"
Environment="NODE_LOG_DIR=/var/log/cluster-os"

# Read cluster key if it exists
EnvironmentFile=-/etc/cluster-os/cluster.env

# Restart policy
Restart=always
RestartSec=10s
StartLimitBurst=5
StartLimitInterval=60s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=node-agent

# Security
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/lib/cluster-os /var/log/cluster-os /etc/cluster-os /etc/wireguard /etc/slurm /etc/munge

[Install]
WantedBy=multi-user.target
