# images/ubuntu/files/systemd/cluster-wifi.service
[Unit]
Description=ClusterOS WiFi Setup
# Run AFTER cloud-init completes - cloud-init may block rfkill
After=systemd-modules-load.service systemd-udevd.service systemd-networkd.service cloud-init.target
Wants=systemd-networkd.service
Requires=systemd-modules-load.service
Before=tailscaled.service tailscale-auth.service cluster-autostart.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Give cloud-init and udev time to settle
ExecStartPre=/bin/sleep 3
# Unblock WiFi BEFORE running setup (in case something blocked it)
ExecStartPre=/usr/sbin/rfkill unblock wifi
ExecStart=/usr/local/bin/cluster-wifi-setup
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
# Don't fail boot if WiFi fails
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
