# images/ubuntu/files/systemd/tailscale-auth.service
[Unit]
Description=ClusterOS Tailscale Auto-Authentication
After=network-online.target tailscaled.service cluster-wifi.service
Wants=network-online.target tailscaled.service
# Only run if credentials file exists
ConditionPathExists=/etc/clusteros/tailscale.env
# Remove the state check - we handle re-auth in the script

[Service]
Type=oneshot
EnvironmentFile=/etc/clusteros/tailscale.env
ExecStart=/usr/local/bin/tailscale-auth
RemainAfterExit=yes
# Give network time to stabilize
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
