# images/ubuntu/files/systemd/tailscale-auth.service
[Unit]
Description=ClusterOS Tailscale Auto-Authentication
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service
ConditionPathExists=/etc/clusteros/tailscale.env
ConditionPathExists=!/var/lib/tailscale/tailscaled.state

[Service]
Type=oneshot
EnvironmentFile=/etc/clusteros/tailscale.env
ExecStart=/bin/bash -c '\
    if [ -n "$TAILSCALE_AUTHKEY" ]; then \
        echo "Authenticating Tailscale with auth key..."; \
        tailscale up --authkey="$TAILSCALE_AUTHKEY" --hostname="$(hostname)" --accept-routes --accept-dns=false; \
    elif [ -n "$TAILSCALE_OAUTH_CLIENT_SECRET" ]; then \
        echo "Authenticating Tailscale with OAuth..."; \
        tailscale up --auth-key="$TAILSCALE_OAUTH_CLIENT_SECRET" --hostname="$(hostname)" --accept-routes --accept-dns=false; \
    else \
        echo "No Tailscale credentials found"; \
        exit 1; \
    fi'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
