# images/ubuntu/files/systemd/tailscale-auth.service
[Unit]
Description=ClusterOS Tailscale Auto-Authentication
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service
ConditionPathExists=/etc/clusteros/tailscale.env

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tailscale-auth
RemainAfterExit=yes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
